# systemd/agent.service
[Unit]
Description=ADK Agent Platform
After=network.target

[Service]
# User running the application
User=ubuntu
Group=ubuntu

# Directory where the code is cloned
WorkingDirectory=/home/ubuntu/google-adk-on-bare-metal

# Environment variables
# Option A: Load from file (Recommended)
EnvironmentFile=/home/ubuntu/google-adk-on-bare-metal/.env
# Option B: Set explicitly
# Environment="DATABASE_URL=postgresql://..."
# Environment="OPENROUTER_API_KEY=sk-..."

# Path to 'uv' or 'python' executable
# Assuming 'uv sync' has been run, we use the virtualenv python
ExecStart=/home/ubuntu/google-adk-on-bare-metal/.venv/bin/python -m server

# Restart settings
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
